# ai/roles/vllm/templates/vllm.service.j2
# Managed by Ansible (tenx.ai.vllm).
[Unit]
Description=vLLM OpenAI-compatible inference server ({{ vllm_served_model_name }})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ vllm_user }}
Group={{ vllm_group }}
WorkingDirectory={{ vllm_home }}
EnvironmentFile={{ vllm_env_file }}
ExecStart={{ vllm_venv_dir }}/bin/vllm serve {{ vllm_model }} \
  --host {{ vllm_host }} \
  --port {{ vllm_port }} \
  --served-model-name {{ vllm_served_model_name }} \
  --max-model-len {{ vllm_max_model_len }} \
  --gpu-memory-utilization {{ vllm_gpu_memory_utilization }} \
  --tensor-parallel-size {{ vllm_tensor_parallel_size }}
{%- if vllm_enable_auto_tool_choice %} \
  --enable-auto-tool-choice \
  --tool-call-parser {{ vllm_tool_call_parser }}
{%- endif %}
{%- for arg in vllm_extra_args %} \
  {{ arg }}
{%- endfor %}

# Model download + load easily exceeds systemd's default 90s start timeout.
TimeoutStartSec=infinity
TimeoutStopSec=120
Restart=always
RestartSec=10
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
